# scripts/config/pspi-firmware-update.service
[Unit]
Description=PSPi ATmega firmware updater (runs once late in boot)
After=multi-user.target
ConditionPathExists=__DRIVERS_BASE__/drivers/update_firmware
ConditionPathExists=__DRIVERS_BASE__/drivers/firmware.hex

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__DRIVERS_BASE__
ExecStartPre=/bin/sh -c 'sleep 60'
ExecStart=__DRIVERS_BASE__/drivers/update_firmware __DRIVERS_BASE__/drivers/firmware.hex
ExecStartPost=/sbin/poweroff
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
