# packer/rpi/services/mouse_64.service
[Unit]
Description=PSPi Mouse Service
After=main_64.service

[Service]
ExecStart=mouse_64
Restart=on-failure

[Install]
WantedBy=multi-user.target
