# packer/rpi/services/osd_32.service
[Unit]
Description=PSPi OSD Service
After=main_32.service

[Service]
ExecStart=osd_32
Restart=on-failure

[Install]
WantedBy=multi-user.target
