# drivers/services/gamepad.service
[Unit]
Description=PSPi Gamepad Service
After=main.service

[Service]
ExecStart=gamepad
Restart=on-failure

[Install]
WantedBy=multi-user.target
