# drivers/services/main.service
[Unit]
Description=PSPi Main Service
After=network.target

[Service]
ExecStart=sudo main
Restart=on-failure

[Install]
WantedBy=multi-user.target

