# drivers/services/mouse.service
[Unit]
Description=PSPi Mouse Service
After=main.service

[Service]
ExecStart=sudo mouse
Restart=on-failure

[Install]
WantedBy=multi-user.target

