# drivers/services/osd.service
[Unit]
Description=PSPi OSD Service
After=main.service

[Service]
ExecStart=osd
Restart=on-failure

[Install]
WantedBy=multi-user.target
