# rpi/services/main_64.service
[Unit]
Description=PSPi Main Service
After=network.target

[Service]
ExecStart=main_64
Restart=on-failure

[Install]
WantedBy=multi-user.target
