# rpi/services/mouse_32.service
[Unit]
Description=PSPi Mouse Service
After=main_32.service

[Service]
ExecStart=mouse_32
Restart=on-failure

[Install]
WantedBy=multi-user.target
