# rpi/services/osd_64.service
[Unit]
Description=PSPi OSD Service
After=main_64.service

[Service]
ExecStart=osd_64
Restart=on-failure

[Install]
WantedBy=multi-user.target
